# deploy/picpop.service
# PicPop System Service
# Copy to /etc/systemd/system/picpop.service on the Radxa device
# Then run: systemctl daemon-reload && systemctl enable picpop

[Unit]
Description=PicPop Photo Booth
After=network.target hostapd.service dnsmasq.service
Wants=hostapd.service dnsmasq.service

[Service]
Type=simple
User=picpop
Group=picpop
WorkingDirectory=/opt/picpop
Environment="PATH=/opt/picpop/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PICPOP_ENV=production"

# Start the server
ExecStart=/opt/picpop/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000

# Restart on failure
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=picpop

[Install]
WantedBy=multi-user.target
